# autocorrector-server.service.tmpl
[Unit]
Description=Autocorrector Server
After=syslog.target

[Service]
Type=exec
RemainAfterExit=yes
ExecStart=/usr/local/bin/autocorrector --user {{.User}}

[Install]
WantedBy=default.target